# data/org.freedesktop.impl.portal.desktop.cosmic.service.in
[Unit]
Description=Portal service (COSMIC implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.cosmic
ExecStart=LIBEXECDIR/xdg-desktop-portal-cosmic
